# crontab
0 0/1 * 1/1 * ? * /usr/bin/python /opt/cronScripts/manager.py
0 0/1 * 1/1 * ? * /usr/bin/python /opt/cronScripts/worker.py
