# crontab
SHELL=/bin/sh
PATH=/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin
BASH_ENV=/container.env
*/1 * * * * /bin/bash /opt/cronScripts/cron.sh >> /var/log/cron.log 2>&1
